# com.user.redautomate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.example.redautomate</string>

    <!-- Use the Python interpreter from your virtual environment -->
    <key>ProgramArguments</key>
    <array>
      <string>/Users/yourusername/path/to/your/project/.venv/bin/python3</string>
      <string>/Users/yourusername/path/to/your/project/redautomate.py</string>
    </array>

    <!-- Schedule to run daily at a specific time (e.g., midnight) -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <!-- Run immediately upon load if the scheduled time was missed -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Redirect output to log files -->
    <key>StandardOutPath</key>
    <string>/tmp/redautomate.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/redautomate.err</string>
  </dict>
</plist>
